# telegram-search-bot.service
[Unit]
Description=Telegram Chinese Search Bot
After=network.target

[Service]
Type=simple
User=your_username
WorkingDirectory=/path/to/chinese-search-bot
Environment="PATH=/path/to/chinese-search-bot/venv/bin"
ExecStart=/path/to/chinese-search-bot/venv/bin/python main.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

